# service/autoscaler.service
[Unit]
Description=KVM autoscaler

[Service]
ExecStart=autoscaler -d

[Install]
WantedBy=multi-user.target
