# service/autoscaler.service
[Unit]
Description=PremiScale autoscaler

[Service]
ExecStart=premiscale --daemon --log-stdout

[Install]
WantedBy=multi-user.target
